# crontab
* * * * * /usr/local/bin/node /usr/src/app/index.js > /proc/1/fd/1 2>/proc/1/fd/2

